# deploy/hetzner-runners/systemd/runner-cleanup.service
[Unit]
Description=Runner Fleet Cleanup (Docker + workdirs)
After=docker.service
Wants=docker.service

[Service]
Type=oneshot
WorkingDirectory=/opt/runner-fleet
EnvironmentFile=-/opt/runner-fleet/.env
Environment=RUNNER_DATA_DIR=/opt/runner-fleet/runner-data
ExecStart=/bin/bash /opt/runner-fleet/cleanup_runner_data.sh
